# scheduling/backup-gitlab.service
[Unit]
Description=Backup Gitlab CE

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/bash /home/core/etc/scripts/backup-gitlab.sh
